# centos_common/files/etc/cron.hourly/chef-client
#!/bin/sh
#
## pull chef config

chef-client

logger "chef-client ran at `date`"
